# GPIOrestart.service
[Unit]
Description=Restart GPIOshortcuts after reboot
[Service]
ExecStart=systemctl restart GPIOShortcuts.service

